# oowv-controller.service
[Unit]
Description=OOWV Controller
After=network.target

[Service]
Type=simple
User=<ersetzen>
ExecStart=/home/<ersetzen>/venv/bin/python /home/<ersetzen>/venv/brainTankController/oowvcontroller.py
WorkingDirectory=/home/<ersetzen>/venv/brainTankController
Restart=on-failure

[Install]
WantedBy=multi-user.target