# oowv-controller.service
[Unit]
Description=OOWV Controller

[Service]
Type=simple
ExecStart=/usr/bin/python /home/Smartwatertank-Controller/oowvcontroller.py

[Install]
WantedBy=multi-user.target